# data/etc/systemd/system/coreos-ipxe.service
[Unit]
Description=CoreOS iPXE server

[Service]
ExecStart=/usr/local/bin/coreos-ipxe-server

[Install]
WantedBy=multi-user.target
